# systemd/easyplay-dashboard.service
[Unit]
Description=EasyPlay status dashboard (Flask, local mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=%REPO%
ExecStart=/usr/bin/python3 %REPO%/tools/dashboard/dashboard.py --local
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
