# systemd/nightly-reboot.timer
[Unit]
Description=Trigger nightly-reboot.service every night at 03:00 local time

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=false
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
